# src/core/services/xzWm.service
[Unit]
Description=AI System Cleanup Service
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash -c '\
    journalctl --vacuum-time=2d; \
    find /var/log/cursor -type f -mtime +7 -delete; \
    find /tmp -type f -atime +1 -delete; \
    paccache -rk1; \
    pacman -Sc --noconfirm; \
    rm -f /var/lib/pacman/db.lck; \
    systemctl restart system-ai.slice; \
    sync'
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
MemoryLimit=1G
Slice=system-ai.slice

[Install]
WantedBy=multi-user.target 