# src/0sfI.slice
[Slice]
Description=Project Claude Optimization Slice
# Optimize for focused project work
CPUWeight=85
MemoryHigh=6G
MemoryMax=10G
# Lower IO priority but still responsive
IOWeight=75
IOSchedulingClass=best-effort
IOSchedulingPriority=5

[Install]
WantedBy=cursor.slice 