# src/0tZa.service
[Unit]
Description=Ensure Gauntlet directories exist
After=local-fs.target
Before=gauntlet-optimizer.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'mkdir -p /home/jon/.config/gauntlet /home/jon/.local/state/gauntlet /home/jon/.local/share/gauntlet /home/jon/brainlifts'
User=jon
Group=jon

[Install]
WantedBy=multi-user.target 