# src/37uu.slice
[Unit]
Description=AI Evolution Service Group
Documentation=file:/home/jon/ai_system_evolution/README.md
Before=cursor-context.slice
CollectMode=inactive-or-failed

[Slice]
CPUWeight=10
MemoryHigh=4G
MemoryMax=6G
TasksMax=64
IOWeight=20 