# src/5X5e.service
[Unit]
Description=AI System Autonomic Manager
After=cursor.slice ai-evolution.service
PartOf=cursor.slice
Requires=ai-evolution.service

[Service]
Type=notify
NotifyAccess=all
Environment=PYTHONPATH=/home/jon/ai_system_evolution
Environment=EVOLUTION_ROOT=/home/jon/ai_system_evolution
ExecStart=/usr/bin/python3 -u /home/jon/ai_system_evolution/lib/autonomic/manager.py
Restart=always
RestartSec=10
WatchdogSec=30

# Resource limits
CPUWeight=30
CPUQuota=30%
MemoryHigh=2G
MemoryMax=4G
TasksMax=16
IOWeight=30

[Install]
WantedBy=default.target 