# src/69IY.service
[Unit]
Description=AI-First Development System Optimizer
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/jon/scripts/ai_first_optimizer.sh
User=root
IOSchedulingClass=best-effort
IOSchedulingPriority=0
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
Slice=system-ai.slice

[Install]
WantedBy=multi-user.target 