# src/6I1H.service
[Unit]
Description=Context Optimizer Service
Documentation=file:/home/jon/brain/systems/autonomic.md
After=network.target meta-learner.service essence-harmonizer.service
PartOf=cursor.slice

[Service]
Type=notify
NotifyAccess=all
ExecStart=/home/jon/scripts/cursor/context-optimizer
Restart=always
RestartSec=30
WatchdogSec=300

# Resource management
CPUWeight=35
CPUQuota=25%
MemoryHigh=384M
MemoryMax=768M
TasksMax=24

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.local/share/cursor
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target 