# src/6rKS.service
[Unit]
Description=Meta-Learner Service
Documentation=file:/home/jon/brain/systems/autonomic.md
After=network.target
PartOf=cursor.slice

[Service]
Type=simple
ExecStart=/home/jon/scripts/cursor/meta-learner --sync-knowledge
Restart=always
RestartSec=30

# Resource management
CPUWeight=30
CPUQuota=20%
MemoryHigh=256M
MemoryMax=512M
TasksMax=16

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.local/share/cursor
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target 