# src/7y8F.service
[Unit]
Description=Gauntlet Focus Mode Manager
After=gauntlet-dirs.service
PartOf=gauntlet-focus.slice

[Service]
Type=simple
ExecStart=/usr/local/bin/focus-mode-manager
Slice=gauntlet-focus.slice
Environment=GAUNTLET_CONFIG_DIR=/home/jon/.config/gauntlet
Environment=GAUNTLET_STATE_DIR=/home/jon/.local/state/gauntlet
Environment=GAUNTLET_DATA_DIR=/home/jon/.local/share/gauntlet
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
User=jon
Group=jon

# Resource limits aligned with focus needs
CPUQuota=50%
MemoryHigh=2G
TasksMax=32
IOWeight=60

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.local/state/gauntlet /home/jon/.local/share/gauntlet
PrivateTmp=yes

[Install]
WantedBy=default.target 