# src/8lgq.service
[Unit]
Description=Gauntlet Focus State Manager
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/home/jon/scripts/gauntlet/focus_state_manager.sh monitor
Restart=always
RestartSec=10
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=graphical-session.target 