# src/CllM.service
[Unit]
Description=Gauntlet Focus State Manager
After=cursor.slice meta-learner.service essence-harmonizer.service focus-enhancer.service
PartOf=cursor.slice
Requires=meta-learner.service essence-harmonizer.service

[Service]
Type=notify
NotifyAccess=all
ExecStart=/home/jon/scripts/gauntlet/focus_state_manager.sh monitor
Restart=always
RestartSec=10
WatchdogSec=60

# Resource management
CPUWeight=30
CPUQuota=30%
MemoryHigh=1G
MemoryMax=2G
IOWeight=30
Nice=10

# Integration paths
Environment=CURSOR_CONFIG_DIR=/home/jon/.config/cursor
Environment=CURSOR_DATA_DIR=/home/jon/.local/share/cursor
Environment=GAUNTLET_ROOT=/home/jon/.gauntlet
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.gauntlet/state /home/jon/.gauntlet/metrics /home/jon/.local/share/cursor/contexts
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=false
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true

[Install]
WantedBy=cursor.slice 