# src/CpUr.service
[Unit]
Description=Cursor Semantic Indexer Service
After=cursor-agent.service
Wants=cursor-agent.service
Slice=cursor-context.slice

[Service]
Type=simple
User=jon
Group=jon
Environment=CURSOR_CONTEXT_ROOT=/home/jon/.config/cursor
Environment=CURSOR_DATA_ROOT=/home/jon/.local/share/cursor
ExecStart=/home/jon/scripts/cursor/semantic_indexer.sh
Restart=always
RestartSec=30
Nice=10

[Install]
WantedBy=multi-user.target 