# src/DJos.timer
[Unit]
Description=Cursor Auto-Update Timer
Requires=cursor-updater.service

[Timer]
OnBootSec=5min
OnUnitInactiveSec=1h
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target 