# src/faGH.service
[Unit]
Description=Cursor Auto-Update Service
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/home/jon/scripts/cursor_auto_updater.sh
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
Nice=19
Slice=system-ai.slice

[Install]
WantedBy=multi-user.target 