# src/FBUX.slice
[Unit]
Description=Gauntlet Focus Tools Slice
Before=slices.target
Requires=system-ai.slice
After=system-ai.slice

[Slice]
CPUWeight=90
IOWeight=90
MemoryHigh=8G
TasksMax=128 