# src/fe3G.service
[Unit]
Description=Cursor Codebase Indexing Service
After=network.target

[Service]
Type=oneshot
User=jon
Group=jon
Environment="CURSOR_WORKSPACE=/home/jon/workspace"
ExecStart=/home/jon/scripts/index_codebase.sh
Nice=19
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
MemoryLimit=2G
Slice=system-ai.slice

[Install]
WantedBy=multi-user.target 