# src/fXjG.service
[Unit]
Description=Autonomic Manager Service
Documentation=file:/home/jon/brain/systems/autonomic.md
After=network.target meta-learner.service essence-harmonizer.service
PartOf=cursor.slice

[Service]
Type=simple
ExecStart=/home/jon/scripts/cursor/autonomic-manager
Restart=always
RestartSec=30
TimeoutStartSec=0

# Resource management
CPUWeight=40
CPUQuota=30%
MemoryHigh=384M
MemoryMax=768M
TasksMax=24

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.local/share/cursor /home/jon/.config/cursor
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target 