# src/FXUY.service
[Unit]
Description=AI System Evolution Autonomic Manager
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=jon
Environment=PYTHONPATH=/home/jon/ai_system_evolution
Environment=CONTAINER_NAME=ai_system_evolution
ExecStart=/usr/bin/python3 /home/jon/ai_system_evolution/lib/autonomic/manager.py

# Resource management
CPUWeight=90
CPUQuota=80%
MemoryHigh=2G
MemoryMax=4G
TasksMax=128

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target 