# src/HJkp.service
[Unit]
Description=Gauntlet Learning Integration Service
After=gauntlet-metrics.service
Wants=gauntlet-metrics.service

[Service]
Type=simple
ExecStart=/home/jon/workspace/system-ops/scripts/gauntlet_learner.sh
Restart=always
RestartSec=10
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/jon/.Xauthority

[Install]
WantedBy=default.target 