# src/IMBi.slice
[Unit]
Description=Cursor Autonomic Resource Management
Documentation=https://sre.google/sre-book/table-of-contents/
Before=cursor-context.slice

[Slice]
CPUWeight=95
IOWeight=95
MemoryHigh=12G
MemoryMax=14G
TasksMax=200

# Autonomic Controls
ResourceControl=true
ResourcePolicy=balanced
Delegate=true 