# src/kidR.slice
[Unit]
Description=AI System Slice
Before=slices.target
DefaultDependencies=no
Requires=-.slice
After=-.slice

[Slice]
CPUWeight=100
IOWeight=100
MemoryHigh=infinity
MemoryMax=infinity
TasksMax=infinity
CPUQuota=100%
AllowedCPUs=0-7
MemoryLow=4G
IODeviceWeight=/dev/nvme0n1 500
Delegate=memory pids cpu io 