# src/kXFs.service
[Unit]
Description=Gauntlet Metrics Service
After=gauntlet-optimizer.service
Requires=gauntlet-optimizer.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStartPre=/bin/sh -c 'mkdir -p ${HOME}/.gauntlet/{contexts,metrics,cache,state}'
ExecStart=${HOME}/workspace/system-ops/scripts/gauntlet_metrics.sh

[Install]
WantedBy=default.target 