# src/Lqy0.service
[Unit]
Description=Focus Knowledge Network Integration Service
After=cursor.slice meta-learner.service essence-harmonizer.service gauntlet-focus.service
PartOf=cursor.slice
Requires=meta-learner.service essence-harmonizer.service gauntlet-focus.service

[Service]
Type=notify
NotifyAccess=all
ExecStart=/home/jon/scripts/cursor/focus-knowledge-integrator
Restart=always
RestartSec=10
WatchdogSec=60

# Resource management
CPUWeight=20
CPUQuota=20%
MemoryHigh=1G
MemoryMax=2G
IOWeight=20
Nice=19

# Integration paths
Environment=CURSOR_CONFIG_DIR=/home/jon/.config/cursor
Environment=CURSOR_DATA_DIR=/home/jon/.local/share/cursor
Environment=GAUNTLET_ROOT=/home/jon/.gauntlet
Environment=BRAIN_DIR=/home/jon/brain

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.gauntlet/state /home/jon/.local/share/cursor/contexts /home/jon/.local/share/cursor/crystallized
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=false
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true

[Install]
WantedBy=cursor.slice 