# src/MBen.service
[Unit]
Description=Cursor Editor and Composer Agent Priority Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/jon/scripts/cursor-priority.sh
User=root

[Install]
WantedBy=multi-user.target 