# src/N62d.service
[Unit]
Description=Gauntlet AI System Evolution Autonomic Manager
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=jon
Group=jon
WorkingDirectory=/home/jon/ai_system_evolution
Environment=PYTHONPATH=/home/jon/ai_system_evolution
ExecStart=/usr/bin/python3 -m lib.autonomic.manager
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryHigh=4G
MemoryMax=8G
CPUQuota=200%
TasksMax=100
IOWeight=100

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target 