# src/naGd.service
[Unit]
Description=AI System Evolution Context Service
After=cursor-context.slice
PartOf=ai-evolution.slice

[Service]
Type=notify
WorkingDirectory=/home/jon/ai_system_evolution
ExecStart=/usr/bin/python3 -m services.context_service
Nice=19
CPUQuota=20%
MemoryHigh=1G
IOWeight=10
Restart=always
RestartSec=5
Environment=PYTHONPATH=/home/jon
Environment=CURSOR_CONFIG_DIR=/home/jon/.config/cursor

[Install]
WantedBy=default.target 