# src/o54I.service
[Unit]
Description=Essence System Harmonizer Service
Documentation=file:/home/jon/brain/systems/autonomic.md
After=network.target cursor-semantic-indexer.service cursor-context-optimizer.service
PartOf=cursor.slice

[Service]
Type=notify
NotifyAccess=all
ExecStart=/home/jon/scripts/cursor/essence-harmonizer
Restart=always
RestartSec=30
TimeoutStartSec=0
WatchdogSec=0

# Resource management
CPUWeight=50
CPUQuota=25%
MemoryHigh=512M
MemoryMax=1G
TasksMax=32

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/jon/.local/share/cursor /home/jon/.config/cursor
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=default.target 