# src/OtZg.service
[Unit]
Description=Autonomic Manager for Claude Instances
After=cursor-claude-project.slice cursor-claude-system.slice
Wants=cursor-claude-project.slice cursor-claude-system.slice

[Service]
Type=simple
ExecStart=/home/jon/scripts/cursor/claude-autonomic-manager
Restart=always
RestartSec=5
Environment=CURSOR_CONFIG=/home/jon/.config/cursor
Environment=CURSOR_SHARE=/home/jon/.local/share/cursor

# Resource Management
Nice=-10
CPUWeight=90
IOWeight=90
MemoryHigh=1G

# Protection
ProtectSystem=strict
ReadWritePaths=/home/jon/.local/share/cursor /home/jon/.config/cursor
NoNewPrivileges=yes
RestrictRealtime=yes

[Install]
WantedBy=default.target 