# src/QsZH.service
[Unit]
Description=Cursor Autonomic Manager Service
After=cursor-context.slice
PartOf=cursor-context.slice
Documentation=https://sre.google/sre-book/table-of-contents/

[Service]
Type=simple
ExecStart=%h/scripts/cursor/autonomic-manager
Slice=cursor-context.slice
Environment=CURSOR_CONFIG_DIR=%h/.config/cursor
Environment=CURSOR_STATE_DIR=%h/.local/state/cursor
Environment=CURSOR_DATA_DIR=%h/.local/share/cursor
Environment=ESSENCE_RESOURCES_DIR=%h/.config/cursor/contexts
Restart=always
RestartSec=1
WatchdogSec=60
NotifyAccess=all

# Performance Optimization
CPUQuota=20%
MemoryHigh=2G
MemoryMax=3G
TasksMax=32
IOWeight=40
Nice=10

# Context Preservation
TimeoutStartSec=30
TimeoutStopSec=30
RuntimeMaxSec=86400
RuntimeDirectory=cursor-autonomic
RuntimeDirectoryMode=0750

[Install]
WantedBy=default.target 