# src/RTx3.slice
[Slice]
Description=System Claude Optimization Slice
# Higher CPU priority for filesystem operations
CPUWeight=75
# More memory for caching filesystem data
MemoryHigh=8G
MemoryMax=12G
# Higher IO priority for faster filesystem access
IOWeight=85
# Optimize for background IO
IOSchedulingClass=best-effort
IOSchedulingPriority=4

[Install]
WantedBy=cursor.slice 