# src/ruD0.slice
[Unit]
Documentation=man:systemd.slice(5)
Before=slices.target

[Slice]
# Higher CPU priority for filesystem operations
CPUWeight=75
# More memory for caching filesystem data
MemoryHigh=4G
MemoryMax=6G
# Prevent excessive swapping
MemorySwapMax=2G
# Higher IO priority for faster filesystem access
IOWeight=85

[Install]
WantedBy=cursor.slice 