# src/sXFH.slice
[Unit]
Documentation=man:systemd.slice(5)
Before=slices.target

[Slice]
# Optimize for focused project work
CPUWeight=85
MemoryHigh=4G
MemoryMax=6G
# Prevent excessive swapping
MemorySwapMax=2G
# Lower IO priority but still responsive
IOWeight=75

[Install]
WantedBy=cursor.slice 