# src/tkAQ.service
[Unit]
Description=Cursor AI Agent Service
After=network.target

[Service]
Type=simple
User=jon
Group=jon
Environment="CURSOR_AGENT_PRIVILEGED=1"
Environment="CURSOR_YOLO_MODE=1"
Environment="CURSOR_AGENT_NO_CONFIRM=1"
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RemoveIPC=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
Slice=system-ai.slice

[Install]
WantedBy=default.target 