# src/uCdT.service
[Unit]
Description=Cursor Context Optimizer Service
After=cursor-agent.service cursor-semantic-indexer.service
Wants=cursor-agent.service cursor-semantic-indexer.service
Slice=cursor-context.slice

[Service]
Type=simple
User=jon
Group=jon
Environment=CURSOR_CONTEXT_ROOT=/home/jon/.config/cursor
Environment=CURSOR_DATA_ROOT=/home/jon/.local/share/cursor
ExecStart=/home/jon/scripts/cursor/context_optimizer.sh
Restart=always
RestartSec=30
Nice=15

[Install]
WantedBy=multi-user.target 