# src/UN2j.slice
[Unit]
Description=AI System Slice
Before=slices.target

[Slice]
CPUWeight=100
IOWeight=100
MemoryHigh=infinity
MemoryMax=infinity
TasksMax=infinity
CPUQuota=100% 