# src/W2yj.service
[Unit]
Description=Essence System Integration Service
After=cursor-semantic-indexer.service
PartOf=cursor-context.slice
Documentation=https://sre.google/sre-book/table-of-contents/

[Service]
Type=simple
ExecStart=/usr/local/bin/essence-integrator
Slice=cursor-context.slice
Environment=CURSOR_CONFIG_DIR=/home/jon/.config/cursor
Environment=CURSOR_STATE_DIR=/home/jon/.local/state/cursor
Environment=CURSOR_DATA_DIR=/home/jon/.local/share/cursor
Environment=ESSENCE_RESOURCES_DIR=/home/jon/.config/cursor/contexts
Restart=on-failure
RestartSec=5
User=jon
Group=jon

# Resource Management
CPUQuota=20%
MemoryHigh=2G
TasksMax=16
IOWeight=30
Nice=10

[Install]
WantedBy=default.target 