# src/XwRY.service
[Unit]
Description=Cursor Startup Manager
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/home/jon/scripts/cursor/launch-dual-claude
TimeoutStartSec=60
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target 