# src/Y2yT.slice
[Slice]
Description=Project Claude Optimization Slice
CPUWeight=80
MemoryHigh=6G
IOWeight=90

[Install]
WantedBy=cursor.slice 