# src/y36I.timer
[Unit]
Description=AI-First Development System Optimizer Timer
Requires=ai-first-optimizer.service

[Timer]
OnBootSec=1min
OnUnitInactiveSec=1h
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target 