# src/ysrB.service
[Unit]
Description=AI System Evolution Service
After=cursor.slice
PartOf=cursor.slice

[Service]
Type=notify
NotifyAccess=all
Environment=PYTHONPATH=/home/jon/ai_system_evolution
Environment=EVOLUTION_ROOT=/home/jon/ai_system_evolution
ExecStart=/usr/bin/python3 -u /home/jon/ai_system_evolution/services/evolution_service.py
Restart=always
RestartSec=10
WatchdogSec=30

# Resource limits
CPUWeight=50
CPUQuota=40%
MemoryHigh=4G
MemoryMax=6G
TasksMax=32
IOWeight=50

[Install]
WantedBy=default.target 