# src/Yyvl.service
[Unit]
Description=Cursor Context Crystallization Service
After=cursor-context.slice
PartOf=cursor-context.slice
Documentation=https://sre.google/sre-book/table-of-contents/

[Service]
Type=notify
NotifyAccess=all
Environment=CURSOR_CONFIG_DIR=%h/.config/cursor
Environment=CURSOR_STATE_DIR=%h/.local/state/cursor
Environment=CURSOR_DATA_DIR=%h/.local/share/cursor
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/home/jon/scripts/cursor/context-crystallizer
Restart=always
RestartSec=60
StartLimitIntervalSec=3600
StartLimitBurst=10
TimeoutStartSec=300
TimeoutStopSec=90
WatchdogSec=600

# Resource limits
MemoryHigh=4G
MemoryMax=3G
CPUWeight=100
IOWeight=90

[Install]
WantedBy=default.target 