# src/zg4E.slice
[Unit]
Description=AI System Resources Slice
Before=slices.target
After=-.slice

[Slice]
# Memory Management
MemoryHigh=4G
MemoryMax=6G
MemorySwapMax=2G

# CPU Management
CPUWeight=120
CPUQuota=80%

# IO Management
IOWeight=100

# Task Management
TasksMax=50 