# src/ZRu2.service
[Unit]
Description=Gauntlet AI-Enhanced Optimization Service
After=cursor-context.slice cursor-ai.slice
PartOf=cursor-context.slice

[Service]
Type=notify
NotifyAccess=all
Environment=GAUNTLET_DATA_DIR=%h/.local/share/gauntlet
Environment=CURSOR_CONFIG_DIR=%h/.config/cursor
Environment=CURSOR_STATE_DIR=%h/.local/state/cursor
Environment=CURSOR_DATA_DIR=%h/.local/share/cursor
Environment=AI_MODELS_DIR=%h/.local/share/cursor/models
Environment=CONTEXT_DIR=%h/.local/share/cursor/contexts
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/home/jon/scripts/cursor/gauntlet-optimizer-ai
Restart=always
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=90
WatchdogSec=600

# Resource limits
MemoryHigh=2G
MemoryMax=3G
CPUWeight=100
IOWeight=90

[Install]
WantedBy=default.target 